# src/confs/crontab
*/1 * * * * /bin/bash /scripts/burying_point_oss.sh /data/pc-journal-logs >> /var/log/cron.log 2>&1
